# templates/lupsa.service.j2
[Unit]
Description=Laptop UPS Alert (LUPSA)

[Service]
ExecStart={{ install_path }}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
