# templates/lupsa.service.j2
[Unit]
Description=Laptop UPS Alert (LUPSA)

[Service]
ExecStart={{ install_path }}
StandardOutput=file:{{ log_path }}
StandardError=file:{{ log_path }}

[Install]
WantedBy=multi-user.target
